# xpra/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=devel:devel

[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor
logfile_maxbytes=10MB
logfile_backups=3
loglevel=info

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:xpra]
command=xpra start-desktop %(ENV_DISPLAY)s
    --bind-tcp=0.0.0.0:%(ENV_XPRA_TCP_PORT)s
    --bind-ws=0.0.0.0:%(ENV_XPRA_PORT)s
    --html=on
    --tcp-auth=file:filename=/home/devel/.xpra/password
    --ws-auth=file:filename=/home/devel/.xpra/password
    --start=openbox-session
    --no-daemon
    --no-notifications
    --no-mdns
    --clipboard=yes
    --speaker=yes
    --microphone=no
    --pulseaudio=no
    --dpi=96
    --resize-display=yes
user=devel
autostart=true
autorestart=true
startretries=5
priority=10
stdout_logfile=/var/log/supervisor/xpra-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/xpra-stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=1

[program:supercronic]
command=/usr/local/bin/supercronic-launcher
user=devel
autostart=true
autorestart=true
startretries=5
priority=40
stdout_logfile=/var/log/supervisor/supercronic-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=1
stderr_logfile=/var/log/supervisor/supercronic-stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=1

[program:log-tail]
command=tail -n 0 -F
    /var/log/supervisor/xpra-stdout.log
    /var/log/supervisor/xpra-stderr.log
    /var/log/supervisor/supercronic-stdout.log
    /var/log/supervisor/supercronic-stderr.log
user=devel
autostart=true
autorestart=true
priority=50
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[include]
files = /etc/supervisor/conf.d/*.conf
